# etc/systemd/system/nanowatchdog.service.in
# @configure_input@

[Unit]
Description=NanoWatchdog management daemon
After=local-fs.target network-online.target
Requires=local-fs.target network-online.target

[Service]
Type=forking
PIDFile=/var/run/nanowatchdog.pid
ExecStart=LIBEXECDIR/nw-daemon.pl -verbose
Restart=on-failure

[Install]
WantedBy=multi-user.target
